# files/system/etc/rpm-ostreed-automatic.timer.d/override.conf
[Timer]
OnBootSec=15min
OnUnitInactiveSec=1d
